# Ant-Media-Server-Service/src/main/daemon/antmedia.service
[Unit]
Description=AntMedia

[Service]
Type=simple
ExecStart=/etc/init.d/antmedia

[Install]
WantedBy=multi-user.target
